# addons/supervisor/bustime_turbo_sync.conf
[program:bustime_turbo_sync]
command = /bustime/bustime/.venv/bin/python /bustime/bustime/utils/turbo_sync_adapter.py
user = www-data
autorestart = true
autostart = true
stdout_logfile = /var/log/supervisor/bustime_turbo_sync.log
redirect_stderr = true

