# addons/supervisor/crawlers.conf
[program:c5]
command = /bustime/bustime/.venv/bin/python /bustime/bustime/coroutines/c5.py
user = www-data
autorestart = true
[group:crawlers]
programs=c5
priority=811
