# addons/supervisor/graph_balancer.conf
[program:graph_balancer]
command = /bustime/bustime/.venv/bin/python /bustime/bustime/coroutines/graph_balancer_server.py 50000
user = www-data
autorestart = true