# addons/supervisor/jam.conf
[program:jam]
command = /bustime/bustime/.venv/bin/python /bustime/bustime/coroutines/jam.py
user = www-data
autorestart = false