# addons/supervisor/updaters.conf
[program:updater_5]
command = /bustime/bustime/.venv/bin/python /bustime/bustime/coroutines/updater.py 5
user = www-data
autorestart = true
[group:updaters]
programs=updater_5
priority=790
