# addons/bustime_zwebsocket_pure.conf
[program:bustime_zwebsocket_pure]
command = /mnt/reliable/repos/bustime/bin/python /mnt/reliable/repos/bustime/bustime/zbusd/zwebsocket-pure.py
user = www-data
priority = 900
