# addons/bustime_zproxy.conf
[program:bustime_zproxy]
command = /mnt/reliable/repos/bustime/bin/python /mnt/reliable/repos/bustime/bustime/zbusd/zproxy.py
user = www-data
priority = 910