# addons/bustime_zgate.conf
[program:bustime_zgate]
command = /mnt/reliable/repos/bustime/bin/python /mnt/reliable/repos/bustime/bustime/zbusd/zgate.py
user = www-data
priority = 990