# addons/bustime_rpc_server.conf
[program:bustime_rpc_server]
command = /mnt/reliable/repos/bustime/bin/python /mnt/reliable/repos/bustime/bustime/zbusd/rpc_server.py
user = www-data
priority = 990